# radicale/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# radicale addon
# ==============================================================================
declare -a options

bashio::log.info "Starting nginx in radicale Addon ..."

export INGRESS_URL=$(bashio::addon.ingress_url)

envsubst '$INGRESS_URL' < /etc/nginx/conf.d/default.template > /etc/nginx/conf.d/default.conf || exit 1

nginx -g 'daemon off;'

