# radicale/rootfs/etc/services.d/radicale/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Radicale addon
# ==============================================================================
declare -a options

bashio::log.info "Starting Radicale ..."

mkdir -p /data/radicale/

exec radicale --config /etc/radicale/config 2>&1
