# print-headers/rootfs/etc/services.d/print_headers/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Print headers addon
# ==============================================================================
declare -a options

bashio::log.info "Starting Print headers Addon ..."

export INGRESS_URL=$(bashio::addon.ingress_url)
cd /app
/usr/bin/python3 ./app.py

