# kodi/rootfs/etc/services.d/kodi/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Kodi addon
# ==============================================================================
declare -a options

bashio::log.info "Starting Kodi ..."


export KODI_DATA="/data/kodi"

mkdir -p "${KODI_DATA}"

kodi --windowing=gbm
