# kodi/rootfs/etc/services.d/kodi/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Kodi addon
# ==============================================================================
declare -a options


export KODI_DATA="/data/kodi"
mkdir -p "${KODI_DATA}"


bashio::log.info "Updating settings..."
python3 /etc/services.d/kodi/replace-settings.py

bashio::log.info "Starting Kodi ..."

kodi --windowing=gbm
