# kodi/rootfs/etc/services.d/kodi/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Kodi addon
# ==============================================================================
declare -a options

bashio::log.info "Starting Kodi ..."

kodi-gbm
